# config/crontab/production.crt
LANG=ja_JP.utf8

0 7 * * * cd /var/toire.tv; NODE_ENV=production /usr/local/bin/node bin/backup.js
